# cadvisor/rootfs/etc/s6-overlay/s6-rc.d/cadvisor/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Add-on: cAdvisor
# Runs cAdvisor
# ==============================================================================

bashio::log.info "Starting cAdvisor..."

exec /usr/bin/run_cadvisor.sh
