# promtail/rootfs/etc/services.d/alloy/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Add-on: Alloy
# Runs Grafana Alloy
# ==============================================================================
readonly ALLOY_CONFIG='/data/promtail/config.alloy'
declare log_level

bashio::log.info 'Starting Grafana Alloy...'

journal_path='/var/log/journal'
if ! bashio::fs.directory_exists "${journal_path}" || [ -z "$(ls -A ${journal_path})" ]; then
    bashio::log.info "No journal at ${journal_path}, looking for journal in /run/log/journal instead"
    journal_path='/run/log/journal'
fi

case "$(bashio::config 'log_level')" in \
    trace)      ;& \
    debug)      log_level='debug' ;; \
    notice)     ;& \
    warning)    log_level='warn' ;; \
    error)      ;& \
    fatal)      log_level='error' ;; \
    *)          log_level='info' ;; \
esac;
bashio::log.info "Alloy log level set to ${log_level}"

export "URL=$(bashio::config 'client.url')"
export "JOURNAL_PATH=${journal_path}"
export "LOG_LEVEL=${log_level}"

alloy_args=("run" "${ALLOY_CONFIG}")

# Add server configuration for health checks
alloy_args+=("--server.http.listen-addr=0.0.0.0:9080")

# Set storage path to a writable directory (matches Dockerfile WORKDIR)
alloy_args+=("--storage.path=/data/promtail")

# Note: Alloy doesn't use --log.level flag, logging is controlled via config file
# The LOG_LEVEL environment variable is used in the Alloy configuration

# Ensure config directory and files are accessible
chmod -R 755 /etc/alloy
chmod 644 /etc/alloy/*.alloy

# Ensure data directory exists and is writable
mkdir -p /data/promtail
chmod 755 /data/promtail

# Copy config file to data directory where AppArmor allows full access
cp /etc/alloy/config.alloy "${ALLOY_CONFIG}"
chmod 644 "${ALLOY_CONFIG}"

# Debug: Check current user and file permissions
bashio::log.info "Current user: $(whoami)"
bashio::log.info "Current UID: $(id -u)"
bashio::log.info "Config file permissions:"
ls -la "${ALLOY_CONFIG}"
bashio::log.info "Config directory permissions:"
ls -la /etc/alloy/

# Verify config file is readable
if [[ ! -r "${ALLOY_CONFIG}" ]]; then
    bashio::log.error "Config file ${ALLOY_CONFIG} is not readable!"
    exit 1
fi

bashio::log.info "Handing over control to Grafana Alloy..."
cd /data/promtail
/usr/bin/alloy "${alloy_args[@]}"