# promtail/rootfs/etc/services.d/alloy/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Add-on: Alloy
# Runs Grafana Alloy
# ==============================================================================
readonly ALLOY_CONFIG='/etc/alloy/config.alloy'
declare log_level

bashio::log.info 'Starting Grafana Alloy...'

journal_path='/var/log/journal'
if ! bashio::fs.directory_exists "${journal_path}" || [ -z "$(ls -A ${journal_path})" ]; then
    bashio::log.info "No journal at ${journal_path}, looking for journal in /run/log/journal instead"
    journal_path='/run/log/journal'
fi

case "$(bashio::config 'log_level')" in \
    trace)      ;& \
    debug)      log_level='debug' ;; \
    notice)     ;& \
    warning)    log_level='warn' ;; \
    error)      ;& \
    fatal)      log_level='error' ;; \
    *)          log_level='info' ;; \
esac;
bashio::log.info "Alloy log level set to ${log_level}"

export "URL=$(bashio::config 'client.url')"
export "JOURNAL_PATH=${journal_path}"
export "LOG_LEVEL=${log_level}"

alloy_args=("run" "${ALLOY_CONFIG}")

# Add server configuration for health checks
alloy_args+=("--server.http.listen-addr=0.0.0.0:9080")

if [ "${log_level}" == "debug" ]; then
    bashio::log.debug "Running Alloy in debug mode..."
    alloy_args+=("--log.level=debug")
else
    alloy_args+=("--log.level=${log_level}")
fi

bashio::log.info "Handing over control to Grafana Alloy..."
/usr/bin/alloy "${alloy_args[@]}"